# deployment/roles/ubuntu-bluez/templates/bluetoothd_supervisor.conf
[program:bluetoothd]
command = /usr/local/libexec/bluetooth/bluetoothd --nodetach --plugin=
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=QUIT
user=root
startsecs=3
priority=999
