# deployment/roles/hub-ubuntu/templates/dhcpd_supervisor.conf
[program:dhcpd]
command = /usr/sbin/dhcpd -f -q
autostart=false
autorestart=true
redirect_stderr=true
stopsignal=QUIT
user=root
startsecs=1
